# backend/stun/stund6.service
[Service]
ExecStart=/usr/local/bin/stunserver --family 6
ReadOnlyDirectories=/
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

[Unit]
After=network.target
