# backend/http/geoipupdate.service
[Unit]
Description=GeoIP2 database updates

[Service]
Type=oneshot
ExecStart=/usr/bin/geoipupdate
