# backend/http/geoipupdate.service
[Service]
ExecStart=/usr/bin/geoipupdate
Type=oneshot

[Unit]
Description=GeoIP2 database updates
